# crontab
#run python script every minutes
* * * * * bash -c "sleep $[ ( $RANDOM\%30 ) ]s" ; /usr/local/bin/python /app/crontest.py > /proc/1/fd/1 2>/proc/1/fd/2

