# crontab
SHELL=/bin/bash
#run python script every minutes
* * * * * sleep $(($RANDOM\%30 +1)) ; /usr/local/bin/python /app/crontest.py > /proc/1/fd/1 2>/proc/1/fd/2

