# crontab
*/5 * * * * /usr/bin/sh /example-scheduled-task.sh